# build/package/dm-network.service
[Unit]
Description=Device Model Network Service (go)


[Service]
ExecStartPre=/usr/bin/mkdir -p /var/run/devicemodel
ExecStart=/usr/bin/networkservice unix /var/run/devicemodel/network.sock
ExecReload=/bin/kill -9 $MAINPID
Type=simple
KillMode = process
TimeoutStartSec=0
TimeoutSec=0
RestartSec=5s
Restart=on-failure

[Install]
WantedBy=multi-user.target

